# data/scope-registry.conf.in
description "Scope registry"
author "Michal Hruby <michal.hruby@canonical.com>"

start on started smart-scopes-proxy
stop on desktop-end

respawn
respawn limit 10 60

exec $SNAP/@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_LIBDIR@/@UNITY_SCOPES_LIB@/scoperegistry
