# data/smart-scopes-proxy.conf.in
description "Smart scopes proxy"
author "Marcus Tomlinson <marcus.tomlinson@canonical.com>"

start on scope-ui-starting
stop on desktop-end

respawn
respawn limit 10 60

expect stop

exec $SNAP/@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_LIBDIR@/@UNITY_SCOPES_LIB@/smartscopesproxy upstart
